# com.kms.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kms.weekly</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/yeojooncho/WIP/KMS/.venv/bin/python</string>
        <string>/Users/yeojooncho/WIP/KMS/run.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/yeojooncho/WIP/KMS</string>

    <!-- Monday 08:30 KST -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/yeojooncho/WIP/KMS/logs/launchd-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/yeojooncho/WIP/KMS/logs/launchd-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
